# udev/99-shuttlexpress.rules
ACTION=="add", ATTRS{idVendor}=="0b33", ATTRS{idProduct}=="0020", RUN+="/bin/bash /home/robert/git/ShuttleXpress/run.sh"