# udev/99-shuttlexpress.rules
SUBSYSTEM=="input", ATTRS{idVendor}=="0b33", ATTRS{idProduct}=="0020", MODE="0660", GROUP="shuttle", RUN+="/sbin/modprobe uinput"
